# pitrader.service
[Unit]
Description=PiTrader - Bot de signaux trading
After=network.target ollama.service
Wants=ollama.service

[Service]
Type=simple
User=deer
WorkingDirectory=/home/deer/PiTrader
ExecStart=/usr/bin/python3 /home/deer/PiTrader/main.py --loop --interval 3600
Restart=always
RestartSec=60

# Limites ressources pour Raspberry Pi
CPUQuota=80%
MemoryMax=1G

# Logs
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pitrader

# Variables d'environnement
EnvironmentFile=/home/deer/PiTrader/.env

[Install]
WantedBy=multi-user.target
